# etc/toughee.conf
[unix_http_server]
file=/var/taurusradius/supervisor.sock

[inet_http_server]
port=127.0.0.1:9001        
username=ctlman              
password=ctlroot               

[supervisord]
nodaemon=false
logfile=/var/taurusradius/supervisord.log
logfile_maxbytes=1MB
logfile_backups=32
loglevel=error
pidfile=/var/taurusradius/supervisord.pid
user=taurusradius

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/taurusradius/supervisor.sock

[program:redis]
command=/opt/taurusradius/bin/redis-server --bind 127.0.0.1 --port 16370 --save 60 1 --dir /var/taurusradius --dbfilename toughee.rdb
startretries = 10
autorestart = true

[program:manage]
command=/opt/taurusradius/radiusctl manage -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 5
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/manage.log

[program:ssportal]
command=/opt/taurusradius/radiusctl ssportal -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/ssportal.log


[program:auth]
command=/opt/taurusradius/radiusctl auth -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/authroxy.log

[program:acct]
command=/opt/taurusradius/radiusctl acct -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/acctproxy.log

[program:worker]
command=/opt/taurusradius/radiusctl worker -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
process_name = %(program_name)s%(process_num)d
numprocs=2
redirect_stderr=true
stdout_logfile=/var/taurusradius/worker.log

[program:task]
command=/opt/taurusradius/radiusctl task -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/task.log

[program:syncd]
command=/opt/taurusradius/radiusctl syncd -c /opt/taurusradius/etc/toughee.json
dictionary=/opt/taurusradius
startretries = 10
autorestart = true
redirect_stderr=true
stdout_logfile=/var/taurusradius/syncd.log


[include]
files = conf.d/*.conf
